# scripts/packaging/com.archon.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.archon.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>-m</string>
        <string>src.archon_monitor</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__ARCHON_ROOT__</string>
    <key>StandardOutPath</key>
    <string>__ARCHON_HOME__/logs/monitor-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__ARCHON_HOME__/logs/monitor-stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
