# scripts/packaging/com.archon.benchmark.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.archon.benchmark</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>scripts/benchmark/run-benchmark.py</string>
        <string>--model</string>
        <string>sonnet</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ARCHON_ROOT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__ARCHON_HOME__/logs/benchmark-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__ARCHON_HOME__/logs/benchmark-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>KeepAlive</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
