# scripts/archon/archon-leann-drain.service
[Unit]
Description=Archon LEANN Index Drain

[Service]
Type=oneshot
WorkingDirectory=__ARCHON_HOME__
Environment=HOME=__ARCHON_HOME__
Environment=ARCHON_PROJECT_ROOT=__ARCHON_ROOT__
Environment=PATH=__ARCHON_HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash __ARCHON_HOME__/.archon/scripts/leann-drain.sh
StandardOutput=append:__ARCHON_HOME__/.archon/logs/systemd-leann-drain-stdout.log
StandardError=append:__ARCHON_HOME__/.archon/logs/systemd-leann-drain-stderr.log
Nice=15
