# scripts/archon/archon-rc-prefilter.timer
[Unit]
Description=Archon RocketChat Pre-filter Timer (every 2 min)

[Timer]
OnBootSec=30
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
